# cmi-chat.service
[Unit]
Description=Compiled Modular Intelligence (CMI) Interactive Local Service Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/drawson/deepseek_experiments
ExecStart=/home/drawson/anaconda3/envs/open-webui/bin/python /home/drawson/deepseek_experiments/chat_fused.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cmi-chat-local

# Environment setups
Environment="HOME=/home/drawson"
Environment="PATH=/home/drawson/anaconda3/envs/open-webui/bin:/usr/local/bin:/usr/bin:/bin"

# Resource limits designed for user workspace constraints
MemoryMax=8G
CPUQuota=30%

[Install]
WantedBy=default.target
